# api/docker/crontab
0 * * * * /app/api/manage.py prunematches
